# deploy/cert-check.service
[Unit]
Description=Portfolio TLS-Zertifikats-Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/portfolio-cert-check.sh
