# deploy/cert-check.timer
[Unit]
Description=Täglicher Portfolio-TLS-Zertifikats-Check

[Timer]
OnCalendar=*-*-* 06:30:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
